# templates/hdmi_sleep.service.j2
[Unit]
Description=HDMI Sleep Control
After=network.target

[Service]
ExecStart={{ venv_python }} /usr/local/bin/hdmi_sleep.py
Restart=always
User=pi

[Install]
WantedBy=multi-user.target