# templates/slideshow.service.j2
[Unit]
Description=Slideshow Display
After=network.target

[Service]
ExecStart={{ venv_python }} /usr/local/bin/slideshow.py
Restart=always
User=pi

[Install]
WantedBy=multi-user.target